# .oblivious/services/rest/run
#!/bin/sh
FLASK_APP=rest-server.py flask run --host=0.0.0.0